# etc/s6/nginx/run
#!/bin/bash

if test -f ./setup; then
    source ./setup
fi

exec gosu root /usr/sbin/nginx -g "daemon off;"
